# deploy/valkey.container
# /etc/containers/systemd/valkey.container

[Unit]
Description=Valkey 9.0.1 Key-Value Store
After=network-online.target
Requires=valkey-data.volume

[Container]
# 1. The Image Contract
Image=localhost/valkey:9.0.1

# 2. The User Contract
# We match the UID defined in your build (999)
User=999
Group=999

# 3. The Storage Contract
# Mounts the volume to the container's data directory
Volume=valkey-data.volume:/var/lib/valkey/data

# 4. Networking
PublishPort=6379:6379

# 5. Environment & Config
# Since we baked the config into the image, we don't strictly need env vars
# unless your entrypoint uses them.
Environment=VALKEY_DATA=/var/lib/valkey/data

# 6. Healthcheck (Optional but Recommended)
# Checks if the server is responding to PING
HealthCmd=valkey-cli ping || exit 1
HealthInterval=30s
HealthRetries=3
HealthStartupDelay=10s

[Service]
# Systemd specific restart policy
Restart=always
TimeoutStartSec=60